# deploy/oss/morndraft-oss-pull.service
[Unit]
Description=Pull and atomically deploy the verified MornDraft OSS release
Documentation=https://github.com/WiseWong6/morndraft-oss/tree/main/deploy/oss
Wants=network-online.target
After=network-online.target nginx.service

[Service]
Type=oneshot
User=morndraft-release
Group=morndraft-release
UMask=0027
Environment=MORNDRAFT_STATIC_ROOT=/srv/morndraft-static
Environment=MORNDRAFT_HEALTH_URL=https://morndraft.com/
EnvironmentFile=-/etc/morndraft-oss/pull.conf
LoadCredential=github_token:/etc/morndraft-oss/credentials/github_token
ExecStart=/usr/bin/python3 /usr/local/libexec/morndraft-oss-pull
TimeoutStartSec=5min
TimeoutStopSec=30s

NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ReadWritePaths=/srv/morndraft-static
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
