# deploy/oss/morndraft-oss-pull.timer
[Unit]
Description=Check for a verified MornDraft OSS release every 30 seconds

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=1s
RandomizedDelaySec=3s
Persistent=yes
Unit=morndraft-oss-pull.service

[Install]
WantedBy=timers.target
